# ansible/playbooks/03_backend/templates/upmon-backend.service.j2
[Unit]
Description=Upmon backend — HTTP API
After=network.target

[Service]
Type=simple
WorkingDirectory=%h/app/upmon/backend
ExecStart=%h/.local/bin/uv run upmon-backend
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
