# ansible/playbooks/03_backend/templates/upmon-backend.service.j2
[Unit]
Description=Upmon backend — HTTP API
After=network.target

[Service]
Type=simple
WorkingDirectory={{ app_dir }}/backend
ExecStart=%h/.local/bin/uv run fastapi run --port {{ lookup('pipe', 'rpass port get ' + inventory_hostname + '/upmon-backend') }} src/upmon_backend/main.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
