# ansible/playbooks/05_agent/templates/upmon-agent-cleanup.timer.j2
[Unit]
Description=Daily cleanup of old access log entries

[Timer]
OnCalendar=daily
Persistent=true
RandomizedDelaySec=3600

[Install]
WantedBy=timers.target
